# recipes-nodes/rbuilder-toggle/rbuilder-toggle.cron
*/1 * * * * root /etc/init.d/rbuilder-toggle start
