# recipes-nodes/config-watchdog/config-watchdog.cron
*/1 * * * * root /etc/init.d/config-watchdog start
