# recipes-nodes/rbuilder-status-watchdog/rbuilder-status-watchdog.cron
*/1 * * * * root /etc/init.d/rbuilder-status-watchdog start
